# khal/.config/systemd/user/vdirsyncer.service
[Unit]
Description=vdirsyncer calendar sync

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
