# khal/.config/systemd/user/vdirsyncer.timer
[Unit]
Description=vdirsyncer calendar sync timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
